# rootfs/etc/s6-overlay/s6-rc.d/sdrmap/run
#!/bin/sh
exec /usr/lib/sdrmapfeeder/sdrmapfeeder.sh
